# launchd/com.nanoclaw-unic.rebuild-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.nanoclaw-unic.rebuild-watcher</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/fedor/nanoclaw-unic/scripts/rebuild-watch.sh</string>
	</array>

	<key>WatchPaths</key>
	<array>
		<string>/Users/fedor/nanoclaw-unic/groups/unic-shared-memory/rebuild-trigger.json</string>
	</array>

	<key>StandardOutPath</key>
	<string>/Users/fedor/nanoclaw-unic/logs/rebuild-watcher.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/fedor/nanoclaw-unic/logs/rebuild-watcher.error.log</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/fedor/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>HOME</key>
		<string>/Users/fedor</string>
	</dict>
</dict>
</plist>
